# doc/decide-platform-riskdata.ini
[program:decide-platform-riskdata-api]
directory=/server/decide-platform-riskdata-api
command=php -d swoole.use_shortname=off ./bin/hyperf.php start
user=dev
autostart=true
autorestart=true
startsecs=1
startretries=3
stderr_logfile=/data/stdlog/decide-platform-riskdata-api-error.log
stdout_logfile=/data/stdlog/decide-platform-riskdata-api.log
logfile_maxbytes=200MB
logfile_backups=5
stdout_logfile_maxbytes=200MB
stdout_logfile_backups=5